# chromeos-base/bluetooth-input-fix/files/bluetooth-input-fix.conf
# Copyright (c) 2023 The FydeOS Authors and Alex313031. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "fix bluetoothd crashes when connecting with input devices"
author          "yang@fydeos.io"

start on started udev-boot
task
oom score -100

exec modprobe -q uhid
